# docker/etc/supervisor.d/crond.conf
[program:crond]
autostart = false
command=/usr/sbin/crond -f
