# docker/etc/supervisor.d/rsyslogd.conf
[program:rsyslogd]
autostart = false
command=/usr/sbin/rsyslogd -n
redirect_stderr=true
stdout_syslog=true

