# docker/etc/supervisor.d/wgpinger.conf
[program:wgpinger]
command=/wgpinger.sh
autostart = false
autorestart=true
startsecs=0
redirect_stderr=true
